# supervisor/kafka.conf
[program:kafka]
command=/opt/services/kafka.sh start
stdout_logfile =/var/log/kafka.log
priority=999
autostart=true
autorestart=true
startretries=100
redirect_stderr=true
