# supervisor/kafka.conf
[program:kafka]
command=service kafka start
autostart=true
autorestart=true
