# supervisor/kafka.conf
[program:kafka]
command=/opt/services/kafka start
autostart=true
autorestart=true
