# supervisor/kafka.ini
[program:kafka]
command=bash /opt/services/kafka start
autostart=true
autorestart=true
