# supervisor/zookeeper.ini
[program:zookeeper]
command=/opt/services/zookeeper.sh start
stdout_logfile =/var/log/zookeeper.log
priority=1
autostart=true
autorestart=false
