# supervisor/zookeeper.conf
[program:zookeeper]
command=service zookeeper start
autostart=true
autorestart=true
