# supervisor/zookeeper.conf
[program:zookeeper]
command=bash /opt/services/zookeeper.sh start
stdout_logfile =/var/log/zookeeper.log
priority=1
autostart=true
autorestart=false
