# supervisor/zookeeper.conf
[program:zookeeper]
command=bash /opt/services/zookeeper start
autostart=true
autorestart=true
